# Interfaz_RPF/server/systemd/rpf-context.service
[Unit]
Description=RPF Context Pipeline — genera contexto y actualiza Open WebUI knowledge base
After=network.target

[Service]
Type=oneshot
User=joselozano
Group=joselozano
WorkingDirectory=/home/joselozano/rpf-ejecucion
ExecStart=/usr/bin/python3 /home/joselozano/rpf-ejecucion/generate_context.py
ExecStart=/usr/bin/python3 /home/joselozano/rpf-ejecucion/upload_context.py
EnvironmentFile=/home/joselozano/rpf-ejecucion/.env

StandardOutput=journal
StandardError=journal

SuccessExitStatus=0 1
