# Interfaz_RPF/server/systemd/rpf-sync.service
[Unit]
Description=RPF SharePoint Sync — descarga incremental de datos al servidor
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=joselozano
Group=joselozano
WorkingDirectory=/home/joselozano/rpf-ejecucion
ExecStart=/usr/bin/python3 /home/joselozano/rpf-ejecucion/sp_sync_daemon.py
EnvironmentFile=/home/joselozano/rpf-ejecucion/.env

StandardOutput=journal
StandardError=journal

# Salida 1 (errores parciales) no marca el timer como fallido
SuccessExitStatus=0 1
