# Interfaz_RPF/server/systemd/rpf-watcher.service
[Unit]
Description=RPF File Watcher — detecta archivos nuevos y notifica n8n
After=network.target

[Service]
Type=simple
User=joselozano
Group=joselozano
WorkingDirectory=/home/joselozano/rpf-ejecucion
ExecStart=/bin/bash /home/joselozano/rpf-ejecucion/file_watcher.sh
EnvironmentFile=/home/joselozano/rpf-ejecucion/.env

# Reiniciar automáticamente si el proceso muere
Restart=always
RestartSec=10s

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
